# files/ddbfiles-cleanup.service
[Unit]
Description=Delete old DDB files
After=network.target

[Service]
ExecStart=/usr/local/bin/ddbfiles-cleanup
